# .config/systemd/user/polybar.service
[Unit]
Description=Polybar

[Service]
Type=simple
ExecStart=/usr/bin/polybar Bar -r 

[Install]
WantedBy=user-graphical.target

