# xorg/.config/systemd/user/bspwm.service
[Unit]
Description=Bspwm

[Service]
Type=simple
ExecStart=/usr/bin/bspwm 

[Install]
RequiredBy=user-graphical.target

