# xorg/.config/systemd/user/sxhkd.service
[Unit]
Description=Hotkey listener

[Service]
Type=simple
ExecStart=/usr/bin/sxhkd

[Install]
WantedBy=user-graphical.target

