# xorg/.config/systemd/user/user-graphical.target
[Unit]
Description=User started their X session successfully
Requires=default.target
After=default.target
